# containers/backend/rootfs/etc/cron.d/default
# События Bitrix
* * * * * www-data php -f /var/www/html/api/local/php_interface/cron_events.php >> /var/log/bitrix.log
